# trunk/startup/rc.local
#!/bin/sh
#
# /etc/rc.d/rc.local:  Local system initialization script.
#
# Put any local startup commands in here.  Also, if you have
# anything that needs to be run at shutdown time you can
# make an /etc/rc.d/rc.local_shutdown script and put those
# commands in there.
#
# California PATH DAS3 Vehicle List
#
# altima_silver
# altima_grey
# audi_silver
# audi_red
#
# DAS3 Startup Command
#
# /home/das3/startup/startup.sh vehicle experiment
#
